# etc/cron.d/apticron
# cron entry for apticron

@weekly root if test -x /usr/sbin/apticron; then /usr/sbin/apticron --cron; else true; fi
